# systemd/roxy-core.service
[Unit]
Description=ROXY Core (always-on background service)
After=default.target
# Prevent restart storms
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=%h/.roxy
ExecStart=%h/.roxy/venv/bin/python %h/.roxy/roxy_core.py
# Verify service is actually ready before reporting success
ExecStartPost=%h/.roxy/scripts/wait-for-ready.sh 30
Restart=on-failure
RestartSec=2
# Allow time for startup + readiness check
TimeoutStartSec=45
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-%h/.roxy/.env

[Install]
WantedBy=default.target
