# voice/roxy-voice.service
[Unit]
Description=Roxy Voice Pipeline - Always-on Voice Interface
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/roxy
Environment="PATH=/opt/roxy/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/roxy/venv/bin/python /opt/roxy/voice/pipeline.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target

