# voice/roxy-voice.service
[Unit]
Description=ROXY Voice Pipeline - Complete Voice Interface
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
WorkingDirectory=/home/mark/.roxy
EnvironmentFile=/home/mark/.roxy/.env
Environment="PATH=/home/mark/.roxy/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/mark/.roxy"
ExecStart=/home/mark/.roxy/venv/bin/python /home/mark/.roxy/voice/pipeline.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=roxy-voice

# Resource limits
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=default.target
